# templates/ddns.service.j2
{{ ansible_managed | comment }}

[Unit]
Description= DDNS job

[Service]
Type=oneshot
ExecStart=/bin/bash /sbin/ddns.sh
